# roles/statecheck/templates/flask_supervisor.conf.j2
[program:app]
command=/home/cumulus/StateExporter.py
user=cumulus
autostart=true
autorestart=true
stdout_logfile=/home/cumulus/app.log
stderr_logfile=/home/cumulus/app.err.log
